# hardware-setup/stm_eth_startup.service
[Unit]
Description=Configure ethernet 4 sercos stm

[Service]
Type=oneshot
ExecStart=/usr/bin/stm.sh

[Install]
WantedBy=multi-user.target
